# src/hub/systemd_service/fom-hub.service
[Unit]
Description=The Fuss-o-meter Hub
After=network.target

[Service]
ExecStart=/usr/bin/python3 -u service.py
WorkingDirectory=/home/pi/fom/src/hub
StandardOutput=inherit
StandardError=inherit
Restart=always
User=pi

[Install]
WantedBy=multi-user.target